# systemd/immich-offsite.timer
[Unit]
Description=Taegliches Immich Offsite-Backup

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
