# systemd/immich-backup.service
[Unit]
Description=Immich lokales Backup auf zweite HDD
Requires=docker.service
After=docker.service
# Beide HDDs muessen gemountet sein, bevor das Backup startet.
RequiresMountsFor=/mnt/immich-data /mnt/immich-backup
OnFailure=immich-notify-fail@%n.service

[Service]
Type=oneshot
ExecStart=/opt/immich/scripts/immich-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
