# systemd/immich-offsite.service
[Unit]
Description=Immich Offsite-Backup (rclone) ins Cloud-/Remote-Ziel
Requires=docker.service
After=docker.service immich-backup.service
RequiresMountsFor=/mnt/immich-backup
OnFailure=immich-notify-fail@%n.service

[Service]
Type=oneshot
ExecStart=/opt/immich/scripts/immich-offsite.sh
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
