# systemd/immich-heartbeat.service
[Unit]
Description=Immich Heartbeat und Selbsttest

[Service]
Type=oneshot
ExecStart=/opt/immich/scripts/immich-heartbeat.sh
